# deploy/loja-utilidades.service
[Unit]
Description=Loja de Utilidades Backend API
Documentation=https://github.com/seu-usuario/loja-materiais-utilidades
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/loja-utilidades/backend
Environment=NODE_ENV=production
Environment=PORT=5000
Environment=DB_HOST=localhost
Environment=DB_PORT=5432
Environment=DB_NAME=loja_utilidades_prod
Environment=DB_USER=postgres
Environment=DB_PASSWORD=senha_segura_producao
Environment=JWT_SECRET=seu_jwt_secret_super_seguro_para_producao_2025
Environment=CORS_ORIGIN=https://seudominio.com

ExecStart=/usr/bin/node server.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=loja-utilidades-backend

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/loja-utilidades/backend/uploads /var/www/loja-utilidades/backend/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Health check
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target 